# le/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# 追加定时任务
cat << "EOF" >> /var/spool/cron/crontabs/root
15 1 * * 5 //root/openwrt-backup.sh
*/10 * * * * /root/DenyPwdHack.sh
EOF

# 覆盖软件源
cat << "EOF" > /etc/opkg/distfeeds.conf
src/gz openwrt_core https://mirrors.cloud.tencent.com/lede/snapshots/targets/x86/64/packages
src/gz openwrt_base https://mirrors.cloud.tencent.com/lede/snapshots/packages/x86_64/base
src/gz openwrt_luci https://mirrors.cloud.tencent.com/lede/releases/18.06.9/packages/x86_64/luci
src/gz openwrt_packages https://mirrors.cloud.tencent.com/lede/snapshots/packages/x86_64/packages
src/gz openwrt_routing https://mirrors.cloud.tencent.com/lede/snapshots/packages/x86_64/routing
src/gz openwrt_telephony https://mirrors.cloud.tencent.com/lede/snapshots/packages/x86_64/telephony
EOF

# 覆盖配置文件
wget -N -P /root http://192.168.10.2/openwrt/openwrt-backup.sh
wget -N -P /etc http://192.168.10.2/openwrt/AdGuardHome.yaml
wget -N -P /etc http://192.168.10.2/openwrt/dnsmasq.conf
wget -N -P /etc http://192.168.10.2/openwrt/hosts
wget -N -P /etc/config http://192.168.10.2/openwrt/config/passwall
wget -N -P /etc/config http://192.168.10.2/openwrt/config/shadowsocksr
#wget -N -P /etc/mosdns http://192.168.10.2/openwrt/mosdns/config.yaml
wget -N -P /etc/mosdns/rule http://192.168.10.2/openwrt/mosdns/hosts.txt

sleep 1

# 重启服务
service cron restart
/etc/init.d/AdGuardHome restart
/etc/init.d/dnsmasq restart
/etc/init.d/shadowsocksr restart
/etc/init.d/passwall restart
mosdns service restart

opkg install http://192.168.10.2/ipk/cdnspeedtest.ipk

opkg install http://192.168.10.2/ipk/luci-app-cloudflarespeedtest_all.ipk

sleep 1

wget -N -P /etc/config http://192.168.10.2/openwrt/config/cloudflarespeedtest

# 清空rc.local
cat << "EOF" > /etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

exit 0
EOF

exit 0
